# one-api.service
# File path: /etc/systemd/system/uniapi.service
# sudo systemctl daemon-reload
# sudo systemctl start uniapi
# sudo systemctl enable uniapi
# sudo systemctl status uniapi
[Unit]
Description=UniAPI Service
After=network.target

[Service]
User=ubuntu  # Make sure to change the username
WorkingDirectory=/path/to/uniapi  # Make sure to change the path
ExecStart=/path/to/uniapi/uniapi --port 3000 --log-dir /path/to/uniapi/logs  # Make sure to change the path and port
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
